# files/etc/uci-defaults/02_network
#!/bin/sh
# Copyright (c) 2015 The Linux Foundation. All rights reserved.
# Copyright (c) 2011-2015 OpenWrt.org

. /lib/functions/uci-defaults.sh
. /lib/functions/system.sh

an7581_setup_interfaces()
{
    local board="$1"

    case "$board" in
    airoha,an7581-evb)
        ucidef_set_interfaces_lan_wan "lan1 lan2 lan3 lan4" "eth1"
        ;;
    bell,xg-040g-md)
        ucidef_set_interfaces_lan_wan "lan2 lan3 lan4" "eth1"
        ;;
    *)
        echo "Unsupported hardware. Network interfaces not initialized"
        ;;
    esac
}

board_config_update
board=$(board_name)
an7581_setup_interfaces $board

# ============================================
# 自定义配置：设置固定 IP 192.168.2.1
# ============================================
uci set network.lan.ipaddr='192.168.2.1'
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.ip6assign='60'

# 可选：设置网关和 DNS（如果需要）
# uci set network.lan.gateway='192.168.2.1'
# uci set network.lan.dns='114.114.114.114'

uci commit network

board_config_flush

exit 0
