# files/etc/config/hotplug.d/USB/10-usb-led
#!/bin/sh
# USB 热插拔脚本 - 控制专用 USB LED（简化版）

[ "$SUBSYSTEM" = "usb" ] || exit 0

case "$ACTION" in
add)
    # USB 插入：两个 USB LED 都点亮
    echo "default-on" > /sys/class/leds/green:usb-1/trigger 2>/dev/null
    echo "default-on" > /sys/class/leds/green:usb-2/trigger 2>/dev/null
    ;;
remove)
    # USB 拔出：两个 USB LED 都熄灭
    echo "none" > /sys/class/leds/green:usb-1/trigger 2>/dev/null
    echo "none" > /sys/class/leds/green:usb-2/trigger 2>/dev/null
    ;;
esac

exit 0
